# salt/application/humanitybox.supervisor.conf
[program:humanitybox]
directory=/var/www/humanitybox
user={{ pillar.deploy_user }}
command={% if pillar.env_name == 'vagrant' %}supervisor -i .git,public,node_modules,salt,scripts{% else %}node{% endif %} server
autostart=true
autorestart=true
{% if pillar.env_name != 'vagrant' %}environment=NODE_ENV=production{% endif %}
